# config/dev/Azure/cordaNotary.service
Description=Corda Node - Notary
Requires=network.target

[Service]
Type=simple
User=acostarodrigo
WorkingDirectory=/home/acostarodrigo/ShieldCorda/build/nodes/Notary
ExecStart=/usr/bin/java -jar /home/acostarodrigo/ShieldCorda/build/nodes/Notary/corda.jar
Restart=on-failure
Environment="CAPSULE_CACHE_DIR=./capsule"

[Install]
WantedBy=multi-user.target
